# crontab
# min hr mday month wday command 
0	*/2	*	*	*	~/.cron.d/areena-rsync.sh
