# crontab
# min hr mday month wday command 
0	*/2	*	*	*	/Users/kkainu/.cron.d/areena-rsync.sh
