# crontab
# min hr mday month wday command 
*	*	*	*	*	/Users/kkainu/.dotfiles/cron.d/areena-rsync.sh
